# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log 
pidfile=/var/run/supervisord.pid ; (supervisord 
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)
user=root

nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[include]
files = /etc/supervisor/conf.d/*.conf

[program:pgadapter]
command=java -jar pgadapter.jar -p %(ENV_PROJECT_ID)s -i %(ENV_INSTANCE_NAME)s -d %(ENV_DATABASE_NAME)s
prioriy=10
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:main]
command=/main
prioriy=20
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
